# service/systemd/spacenav-ws.service
[Unit]
Description=SpaceNav WebSocket bridge for OnShape
After=network.target

[Service]
ExecStart=/usr/bin/env uv run --with spacenav-ws %h/.local/bin/linapse-ws-proxy
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
